# scripts/catscan-deploy.service
[Unit]
Description=CatScan auto-deploy poll (git pull + rebuild + restart on new commit)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=catscan
WorkingDirectory=/opt/catscan
RuntimeDirectory=catscan-deploy
ExecStart=/opt/catscan/scripts/auto-deploy.sh
